# rss.service
[Unit]
Description=RssService

[Service]
ExecStart=/usr/bin/python /home/kouta/Scripts/rss.py

