# wall_changer.service
[Unit]
Description=Wallpaper changer service

[Service]
ExecStart=/bin/bash 'caminho para o script .sh'
